# config/systemd/user/greenclip.service
[Unit]
Description=Clipboard Manager
After=default.target
Requires=default.target
[Service]
Environment="DISPLAY=:0"
ExecStart=/usr/local/bin/greenclip daemon
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
